# deploy/supervisor.conf
[program:mm]
command=/var/www/mm/venv/bin/gunicorn -c deploy/gunicorn.conf wsgi:app
directory=/var/www/mm
user=root
autostart=true
autorestart=true
environment = MODE="PRODUCTION"